# supervisor.d/opendkim.ini
[program:opendkim]
user=opendkim
command=/usr/sbin/opendkim -x /run/opendkim/opendkim.conf
autorestart=true
startsecs=0
