# supervisor.d/opendkim.ini
[program:opendkim]
user=opendkim
command=/usr/sbin/opendkim -x /etc/opendkim/opendkim.conf
autorestart=false
startsecs=0

