# supervisor.d/syslogd.ini
[program:syslog]
command=/sbin/syslogd
autorestart=false
startsecs=0
