# supervisor.d/syslogd.ini
[program:syslog]
command=/sbin/syslogd
autorestart=true
startsecs=0
